# java/padilla/server/simple/src/main/install/systemd/ocpp-server.service
[Unit]
Description=OCPP Server

[Service]
Type=forking
ExecStart=/usr/local/ocpp-server/install/systemd/launch start
ExecStop=/bin/kill -15 $MAINPID

[Install]
WantedBy=reboot.target
