# java/padilla/server/simple/src/main/install/systemd/ocpp-server.service
[Unit]
Description=OCPP Server

[Service]
Type=forking
ExecStart=/usr/local/ocpp-server/install/systemd/saltogateway start
ExecStop=/usr/local/ocpp-server/install/systemd/saltogateway stop

[Install]
WantedBy=reboot.target
