# java/padilla/server/simple/src/main/install/systemd/ocpp-server.service
[Unit]
Description=OCPP Server

[Service]
User=root
Type=forking
ExecStart=/usr/local/ocpp-server/systemd/launch.sh start
ExecStop=/bin/kill -15 $MAINPID

[Install]
WantedBy=reboot.target
